# jail/etc/inittab
tty1::askfirst:/bin/login
:1:respawn:yes>/dev/null
